# data/ubuntu-hibernate-guard.timer
[Unit]
Description=Run hibernate configuration check shortly after boot

[Timer]
OnBootSec=3min
Persistent=false

[Install]
WantedBy=timers.target
