# packaging/postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ]; then
    # Keep hibernation configuration intact.  v0.42 managed config files should
    # be changed only through the wizard rollback/removal flow, not silently by
    # package purge.
    rm -rf /var/lib/ubuntu-hibernate-wizard
    systemctl disable ubuntu-hibernate-guard.timer 2>/dev/null || true
fi
exit 0
